# web/osmx-web.service
[Unit]
Description=osmx-web
After=network-online.target

[Service]
ExecStart=/home/osmx/main -resultsDir=/mnt/www/files -tmpDir=/mnt/tmp -osmxExec=/home/osmx/OSMExpress/osmx -osmxData=/mnt/planet.osmx
Restart=always
User=osmx

[Install]
WantedBy=multi-user.target